# deployment/systemd/ai-trader.service
[Unit]
Description=AI OKX Trader Service
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=trader
WorkingDirectory=/opt/ai-okx-trader
ExecStart=/usr/bin/python3 -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
